# systemd/gpu-switcher.service
[Unit]
Description=Automatic GPU switching service

[Service]
ExecStart=/usr/bin/gpu-switcher

[Install]
WantedBy=graphical.target